# django_bootstrapper/project/{{cookiecutter.directory_name}}/nginx/supervisor.conf
[program:{{cookiecutter.project_slug}}]
command = /webapps/{{cookiecutter.directory_name}}/gunicorn_start.sh
user = {{cookiecutter.user_name}}
stdout_logfile = /webapps/{{cookiecutter.directory_name}}/logs/gunicorn_supervisor.log
redirect_stderr = true
